# config/systemd/rotter@raar.service.d/override.conf
[Service]
# Rotter startup options
Environment="ROTTER_OPTIONS=-a -f flac -c 2 -R 10 -L '%%FT%%H%%M%%S%%z_060.flac' -j"
