# config/systemd/rotter@raar.service.d/override.conf
[Service]
# Rotter startup options
Environment="ROTTER_OPTIONS=-a -f flac -c 2 -L '%%FT%%H%%M%%S%%z_60.flac' -j"
